# www-apps/gitlab-pages/files/gitlab-pages.initd
#!/sbin/openrc-run
# Copyright 2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/opt/gitlab/gitlab-pages/gitlab-pages"
description="Gitlab Pages"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-config /etc/conf.d/gitlab-pages"
start_stop_daemon_args="--background \
	--stdout /var/log/${RC_SVCNAME}.log \
	--stderr /var/log/${RC_SVCNAME}.log"

depend() {
    use net gitlab-puma
}
